# etc/systemd/system/gps-monitor.service
[Unit]
Description=Monitor GPS and log location information
After=network.target
Wants=gpsd.service gpsd.socket

[Service]
ExecStart=/home/pi/gps_monitor.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
