# etc/systemd/system/gps-init.service
[Unit]
Description=Configure u-blox USB GPS and start gpsd
After=network.target

[Service]
ExecStart=/home/pi/setup-gps.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target